# backend/ops/ec2/splice-step4-worker.service
[Unit]
Description=splice-playground STEP4 ColabFold worker
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/splice-playground/backend
EnvironmentFile=/home/ec2-user/splice-playground/backend/.env.worker
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ec2-user/.local/bin/uv run python -u scripts/run_step4_jobs.py --provider colabfold --poll-seconds 15 --limit 2
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
